# remote/template-run.conf
[program:${COMPONENT}]
environment=HOME=/home/easydeploy/
priority= 50
command = /home/easydeploy/bin/run-docker.sh %(process_num)02d
process_name = %(program_name)s-%(process_num)02d
numprocs = ${EASYDEPLOY_PROCESS_NUMBER}
directory=/home/easydeploy/
autostart=false
autorestart=true
startsecs=0
startretries=60
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
user=easydeploy
redirect_stderr=true
stdout_logfile=/var/log/easydeploy/run-%(process_num)02d.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
serverurl=AUTO


[program:datadog_agent]
priority= 40
command = /home/easydeploy/bin/datadog-agent.sh
process_name = %(program_name)s
numprocs = 1
directory=/root
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0
stopsignal=INT
stopwaitsecs=60
stopasgroup=true
killasgroup=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/easydeploy/datadog-agent.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB



[program:discovery]
priority= 20
command =  /home/easydeploy/bin/discovery.sh
process_name = %(program_name)s
numprocs = 1
autostart=true
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/easydeploy/discovery.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
serverurl=AUTO

[program:stats]
priority= 30
command =  /home/easydeploy/bin/stats.sh
process_name = %(program_name)s
numprocs = 1
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/easydeploy/stats.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
serverurl=AUTO

[program:clean_on_boot]
priority= 0
command =  /home/easydeploy/bin/clean.sh
process_name = %(program_name)s
numprocs = 1
autostart=true
autorestart=unexpected
startsecs=0
startretries=10
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/easydeploy/clean_on_boot.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
serverurl=AUTO

[program:serf]
priority= 10
command = /home/easydeploy/bin/serf-agent.sh
process_name = %(program_name)s
numprocs = 1
autostart=true
autorestart=true
startsecs=10
startretries=60
exitcodes=0
stopasgroup=true
stopsignal=TERM
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/easydeploy/serf.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
serverurl=AUTO


[program:consul]
priority= 20
command = /home/easydeploy/bin/consul-agent.sh
process_name = %(program_name)s
numprocs = 1
autostart=true
autorestart=true
startsecs=10
startretries=60
exitcodes=0
stopasgroup=true
stopsignal=TERM
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/easydeploy/consul.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
serverurl=AUTO

[program:weave]
priority= 30
command = /ezbin/weave-agent.sh
process_name = %(program_name)s
numprocs = 1
autostart=true
autorestart=true
startsecs=10
startretries=6000
exitcodes=0
stopasgroup=true
stopsignal=TERM
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/easydeploy/weave.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
serverurl=AUTO

[program:logstash-ship]
priority= 30
command =  /ezbin/logstash-ship.sh
environment=HOME=/root
process_name = %(program_name)s
numprocs = 1
autostart=false
autorestart=true
startsecs=10
startretries=60
exitcodes=0
stopsignal=INT
stopasgroup=true
stopwaitsecs=30
user=root
redirect_stderr=true
stdout_logfile=/var/log/easydeploy/logstash-ship.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
serverurl=AUTO

[program:health-check-server]
priority= 100
command =  /ezbin/health_check_server.sh
environment=HOME=/root
process_name = %(program_name)s
numprocs = 1
autostart=true
autorestart=true
startsecs=0
startretries=100000
exitcodes=0
stopsignal=TERM
stopasgroup=true
stopwaitsecs=30
user=root
redirect_stderr=true
stdout_logfile=/ezlog/health-check-server.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
serverurl=AUTO


[program:btsync]
priority= 100
command = /ezbin/btsync.sh
process_name = %(program_name)s
numprocs = 1
autostart=true
autorestart=true
startsecs=10
startretries=5
exitcodes=0
stopsignal=TERM
stopasgroup=true
stopwaitsecs=30
user=easydeploy
redirect_stderr=true
stdout_logfile=/ezlog/btsync.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
